# vnc/x11vnc.conf
[program:x11vnc]
#command=x11vnc -forever -shared
command=x11vnc -rfbauth /home/app/.vnc/passwd -forever -shared
autorestart=true